# services/rag-pipeline/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:fastapi]
command=uvicorn app.main:app --host 0.0.0.0 --port 8002
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONPATH=/app

[program:llama_server]
command=bash -c 'MODEL_FILE="/models/mistral-7b-instruct-v0.3.Q5_K_M.gguf"; while [ ! -f "$MODEL_FILE" ]; do echo "Waiting for $MODEL_FILE to be mounted..."; sleep 10; done; exec llama-server --model "$MODEL_FILE" --port 8080 --ctx-size 2048 --parallel 1 --n-gpu-layers 0'
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
